# root/etc/uci-defaults/40_luci-webrtc-control
#!/bin/sh
/etc/init.d/webrtc-control enable
/etc/init.d/webrtc-control start
rm -f /tmp/luci-indexcache
exit 0
